# scheduling/com.dealintel.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dealintel.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/deal-intel/.venv/bin/dealintel</string>
        <string>weekly</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/deal-intel</string>

    <!-- Schedule: Run weekly on Sunday at 11:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/path/to/deal-intel/logs/weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/deal-intel/logs/weekly.err</string>

    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
